# private_dot_local/bin/systemd/dropbadips.service
[Unit]
Description="Run dropbadips script"

[Service]
ExecStart=/usr/local/bin/dropbadips
