# private_dot_local/bin/systemd/dropbadips-boot.service
[Unit]
Description="Run dropbadips reset for boot, via associated timer"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dropbadips reset ru cn ir
